# puppet/local_modules/log_generator/templates/apache_access_log_generator.service.erb
[Unit]
Description=Apache access log generator
After=syslog.target
After=network.target

[Service]
Type=simple
Restart=on-failure
User=log_generator
Group=log_generator
WorkingDirectory=<%= @install_dir %>
ExecStart=/bin/bash <%= @install_dir %>/apache_access_log_generator_tcp.sh <%= @flume_agent_host %> <%= @flume_agent_port %> <%= @rows_in_batch %> <%= @sleep_from %> <%= @sleep_to %>

# Give a reasonable amount of time for the server to start up/shut down
TimeoutSec=300

[Install]
WantedBy=multi-user.target